# bin/com.nickorefice.turo-tiguans.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nickorefice.turo-tiguans</string>

    <!--
      Repo lives in ~/code/ (not ~/Documents/) so TCC doesn't gate access.
      launchd execs the script directly — no AppleScript-applet wrapper
      needed in this location. (build-applet.sh remains in bin/ if you
      ever move the repo back into a TCC-protected folder.)
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/nickorefice/code/turo-ai-scaper/bin/run-scrape.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nickorefice/code/turo-ai-scaper</string>

    <!-- Fire once a day at 09:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- If the Mac was asleep at 09:00, run as soon as it wakes -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/nickorefice/code/turo-ai-scaper/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nickorefice/code/turo-ai-scaper/logs/launchd.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>AUTO_COMMIT</key>
        <string>1</string>
    </dict>
</dict>
</plist>
